# minidsp_2x4hd_usb/services.d/minidsp/run
#!/usr/bin/with-contenv bash
# Скрипт запускается под s6-overlay

# Путь к Python скрипту
python3 /app/minidsp_status.py
