# minidsp_2x4hd_usb/services.d/minidsp/run
#!/usr/bin/with-contenv bash
# Используем python напрямую из venv, без активации
exec /venv/bin/python /app/minidsp_status.py
