# minidsp_2x4hd_usb/services.d/minidsp/run
#!/usr/bin/with-contenv bash

# активируем venv
. /venv/bin/activate

# запускаем основной скрипт (exec важно для корректного PID внутри s6)
exec python /app/minidsp_status.py
