# minidsp_2x4hd_usb/services.d/minidsp/run
#!/usr/bin/with-contenv bash

# Запуск основного скрипта
exec /venv/bin/python /app/minidsp_status.py
